# btrfs-scrub.service
[Service]
ExecStart=/usr/local/bin/btrfs-scrub.sh
